# provisioning/roles/redis_data/files/etc/init.d/redis-data.conf
description "redis server"

start on runlevel [23]
stop on shutdown

exec sudo /usr/local/bin/redis-server /etc/redis/data.conf

respawn
